# archez-installer/files/usr/lib/archez-installer/config-first-boot.service
[Unit]
Description=Prepare system for first boot

[Service]
Type=oneshot
ExecStart=/config-first-boot.sh

[Install]
WantedBy=multi-user.target
